# conf/systemd/meta-open-db.com-manager.path
[Path]
PathModified=/var/www/meta-open-database

[Install]
WantedBy=multi-user.target
